# systemd/ntc-cvav-dn700r-recorder-pipeline.timer
[Unit]
Description=Run CVAV portable-rack DN700R recorder pipeline every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=5s
Unit=ntc-cvav-dn700r-recorder-pipeline.service

[Install]
WantedBy=timers.target
