# systemd/ntc-dn700r-recorder-pipeline.service.d/10-run-only-when-idle.conf
[Service]
ExecStart=
ExecStart=/usr/bin/python3 /root/NTC-Runtime/autosyncmix/bin/run_dn700r_pipeline_when_idle.py --label NTC-DN700R --config /root/NTC-AutoSyncMix/config/recorders/DN700R_ftp.json -- /root/NTC-AutoSyncMix/scripts/recorders/run_DN700R_pipeline.sh
ExecStartPost=
