# systemd/ntc-mixassist.service.d/source-label.conf
[Service]
Environment="MIXASSIST_UDP_SOURCE_LABEL=Current WebCall Room LR"
EnvironmentFile=-/run/ntc-mixassist/app-source.env
