# systemd/ntc-mixassist.service.d/source-label.conf
[Service]
Environment="MIXASSIST_UDP_SOURCE_LABEL=Q-SYS NTC-LR"
